# examples/archlinux/etc/watchman/init.d/rc.local
#!/bin/env watchman

service_type='conf'
service_command="$BASH"
service_args='/etc/rc.local'

unset stop

start() {
	watchman.start && watchman.done
}
